# docs/systemd/system/goauthing@.service
[Unit]
Description=Authenticating utility for net.cuc.edu.cn
StartLimitIntervalSec=0

[Service]
# default config is in ~/.auth-cuc
ExecStartPre=-/usr/local/bin/auth-cuc -D deauth
ExecStartPre=-/usr/local/bin/auth-cuc -D auth
ExecStart=/usr/local/bin/auth-cuc online
User=%i
Restart=always
RestartSec=5

[Install]
WantedBy = multi-user.target
